# myprojectmanager-backup.timer
#!/usr/bin/env bash
# systemd timer unit for automatic daily data backup
# Triggers myprojectmanager-backup.service at midnight (00:00) every day

[Unit]
Description=MyProjectManager Daily Data Backup Timer (runs at midnight)
Requires=myprojectmanager-backup.service

[Timer]
# Run daily at midnight (00:00) server local time
OnCalendar=daily
# Alternative: Run at specific time, e.g., 02:30 AM
# OnCalendar=*-*-* 02:30:00

# If the system was off at trigger time, run the backup immediately on next boot
Persistent=true

# Add random delay up to 30 minutes to avoid load spikes if multiple timers run at midnight
RandomizedDelaySec=30m

# Accuracy: Timer may drift within 1 minute (reduces CPU wakeups)
AccuracySec=1min

[Install]
WantedBy=timers.target
